# supervisor/programs.conf
[program:ollama]
command=/usr/local/bin/ollama serve
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/ollama.log
stderr_logfile=/var/log/supervisor/ollama.log
environment=OLLAMA_HOST="%(ENV_OLLAMA_HOST)s",OLLAMA_KEEP_ALIVE="%(ENV_OLLAMA_KEEP_ALIVE)s",OLLAMA_NUM_PARALLEL="%(ENV_OLLAMA_NUM_PARALLEL)s"
priority=10

[program:chromadb]
command=chroma run --host %(ENV_CHROMA_HOST)s --port %(ENV_CHROMA_PORT)s --path %(ENV_CHROMA_DATA_DIR)s
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/chromadb.log
stderr_logfile=/var/log/supervisor/chromadb.log
priority=20

[program:mcp]
command=python3 /scripts/mcp_server.py
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/mcp.log
stderr_logfile=/var/log/supervisor/mcp.log
priority=30

[program:openwebui]
command=open-webui serve --host 0.0.0.0 --port 8080
autostart=true
autorestart=true
priority=40
stdout_logfile=/var/log/openwebui.log
stderr_logfile=/var/log/openwebui.err
environment=OLLAMA_BASE_URL="http://localhost:11434",WEBUI_AUTH=False
